# rootfs/etc/s6-overlay/s6-rc.d/nginx-prepare/up
/etc/s6-overlay/scripts.d/00-nginx-init.sh
