# rootfs/etc/s6-overlay/s6-rc.d/nginx-log/run
#!/bin/sh
exec logutil-service /run/nginx/log
