# rootfs/etc/s6-overlay/s6-rc.d/docker-entrypoint/up
/etc/s6-overlay/scripts.d/docker-entrypoint.sh
