# rootfs/etc/s6-overlay/s6-rc.d/nginx-log-prepare/up
if { mkdir -p /var/log/nginx/persistent }
if { chown nobody:nogroup /var/log/nginx/persistent }
chmod 02755 /var/log/nginx/persistent
