# rootfs/etc/s6-overlay/s6-rc.d/nginx-log-prepare/up
if { mkdir -p /run/nginx/log }
if { chown nobody:nogroup /run/nginx/log }
chmod 02755 /run/nginx/log
